# dirlock/dirlock-sddm.service
[Unit]
Description=Enable the PAM dirlock modules before starting SDDM
Before=sddm.service
PartOf=sddm.service

[Service]
Type=oneshot
ExecStartPre=rm -f /etc/sddm.conf.d/zz-steamos-dirlock.conf
ExecStart=/usr/lib/steamos/dirlock-pam-config

[Install]
WantedBy=sddm.service
